# docker/mail/s6-rc/services/policydRateLimit/run
#!/command/execlineb -P
policyd-rate-limit